# deploy/goblinSlop.service
[Unit]
Description=GoblinSlop - A chaotic collection of goblin knowledge
After=network.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=/home/__APP_USER__
ExecStart=/home/__APP_USER__/goblin_slop
Restart=always
RestartSec=5
Environment="RUST_LOG=info"
Environment="RUST_BACKTRACE=1"
Environment="GOBLIN_HOST=127.0.0.1"
Environment="GOBLIN_PORT=3000"
Environment="GOBLIN_DB_PATH=/home/__APP_USER__/goblin_slop.db"
Environment="GOBLIN_CONTENT_DIR=/home/__APP_USER__/data/content"
Environment="GOBLIN_STATIC_DIR=/home/__APP_USER__/static"
Environment="GOBLIN_DATA_DIR=/home/__APP_USER__/data"

[Install]
WantedBy=multi-user.target